# static/linux/kdr.conf
description "KodeDrive Daemon"

start on startup

task
export PATH=~/.local/bin:"$PATH"
kdr sys -i
